# automation/systemd/worldcup-2026-update.timer
[Unit]
Description=Run World Cup 2026 update daily at 15:00 Asia/Shanghai

[Timer]
OnCalendar=*-*-* 15:00:00 Asia/Shanghai
Persistent=true
Unit=worldcup-2026-update.service

[Install]
WantedBy=timers.target
